# dist/supervisord.conf
[supervisord]
logfile=/opt/isl_client/logs/supervisord.log ; основной лог-файл (по умолчанию $CWD/supervisord.log)
logfile_maxbytes=50MB         ; максимальный размер основного лог-файла перед ротацией (по умолчанию 50MB)
logfile_backups=5             ; количество бэкапов основного лог-файла (по умолчанию 10)
loglevel=debug                ; уровень логирования (по умолчанию info; другие значения: debug, warn, trace)
pidfile=/tmp/supervisord.pid  ; файл PID supervisord (по умолчанию supervisord.pid)
nodaemon=true                 ; запускать в переднем плане, если true (по умолчанию false)
minfds=1024                   ; минимальное количество доступных файловых дескрипторов при запуске (по умолчанию 1024)
minprocs=200                  ; минимальное количество доступных дескрипторов процесса при запуске (по умолчанию 200)

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=admin
password=ats60srb

[program:isl_client]
priority=10
directory=/var/app
command=python runner.py --daemon
user=root
autostart=true
autorestart=true
stderr_logfile=/opt/isl_client/logs/isl_client.err.log
stdout_logfile=/opt/isl_client/logs/isl_client.out.log
