# dist/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700
username = admin
password = ats60srb

[program:isl_client]
command=python3 runner.py --daemon
directory=/opt/isl_client
autostart=true
autorestart=true
stderr_logfile=/opt/isl_client/logs/isl_client.err.log
stdout_logfile=/opt/isl_client/logs/isl_client.out.log
